# filldb-serverslist.crontab
* * * * * /var/www/html/RunMe/filldb-serverslist.php >> /var/log/obi.log 2>&1
